# deploy/etax2-dev-worker.service
[Unit]
Description=etax2 DEV queue worker
After=network.target

[Service]
User=apache
Group=apache
Environment=HOME=/tmp
WorkingDirectory=/var/www/dev.etax2.com
ExecStart=/usr/bin/php artisan queue:work --queue=default --sleep=3 --tries=1 --max-time=3600
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
